# .config/autostart/vboxtray.desktop
[Desktop Entry]
Name=vboxtray
Comment=Virtualbox tray app
Exec=vboxtray %U
Icon=virtualbox.png
Terminal=false
Type=Application
Encoding=UTF-8
