# .config/autostart/vboxtrayico.desktop
[Desktop Entry]
Name=vboxtrayico
Comment=Virtualbox tray app
Exec=vboxtrayico %U
Icon=virtualbox.png
Terminal=false
Type=Application
Encoding=UTF-8
